# freebsd-update-notify.crontab
# Place me in ${PREFIX}/etc/cron.d/freebsd-update-notify

SHELL=/bin/sh
PATH=/etc:/bin:/sbin:/usr/bin:/usr/sbin

00 00 * * * root /usr/local/libexec/freebsd-update-cron > /var/log/freebsd-update-cron 2>&1
00 09,15,21 * * * root /usr/local/libexec/freebsd-update-notify > /var/log/freebsd-update-notify 2>&1
